# provision/service/php-fpm.conf
[program:php-fpm]
command=/usr/sbin/php-fpm7.2 -c /etc/php/7.2/fpm
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
